# systemd/echo.socket
[Unit]
Description=Echo Socket

[Socket]
ListenStream=7
Accept=yes

[Install]
WantedBy=sockets.target
